# etc/supervisor/conf.d/chromium.conf
[program:chromium]
command=/usr/local/bin/chromium-run
priority=99
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
